# roles/docker/files/docker-prune.timer
[Unit]
Description=Docker system prune timer

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
