# roles/auto-update/templates/dnf-autoreboot.timer.j2
[Unit]
Description=upgrade system

[Timer]
OnCalendar={{ auto_update_time }}

[Install]
WantedBy=timers.target
