# roles/auto-update/files/auto-update.service
[Unit]
Description=upgrade system

[Service]
Type=simple
ExecStart=/usr/local/bin/auto-update.sh

