# roles/podman/legacy/nextcloud-cron.service.j2
[Unit]
Description=Nextcloud cron.php job

[Service]
ExecStart=/usr/bin/podman exec nextcloud /usr/local/bin/php -f /var/www/html/cron.php
