# roles/podman/templates/podman-image-prune.timer.j2
[Unit]
Description=podman image prune timer

[Timer]
OnCalendar={{ podman_image_prune_timer }}

[Install]
WantedBy=timers.target

