# roles/podman/templates/podman-system-prune.timer.j2
[Unit]
Description=podman image prune timer

[Timer]
OnCalendar={{ podman_system_prune_timer }}

[Install]
WantedBy=timers.target

