# roles/podman/templates/podman-image-prune.service.j2
[Unit]
Description=prune all unused podman images

[Service]
ExecStart=/usr/bin/podman image prune --all --force

